# templates/default/resque.conf.erb
#!upstart
description "Resque workers"

pre-start script

  mkdir -p /var/log/resque
  chown -R deploy /var/log/resque

end script

start on runlevel [2345]

stop on runlevel [016]
